# systemd/user/auto_theme.service
[Unit]
Description = Automatically change light/dark mode from time of day

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/auto_theme.sh