# systemd/user/auto_theme.timer
[Unit]
Description=Automatically change light/dark mode from time of day

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target